# conf/etc/sv/x/run
#!/bin/sh

exec su -c startx blep >/dev/tty1

