# root/etc/sv/x/run
#!/bin/sh

exec su -c startx x >/dev/tty1

